# server/attendance-server.service
[Unit]
Description=Sistema de Asistencia IoT - FastAPI Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=usuario
Group=usuario
WorkingDirectory=/home/usuario/facial-attendance-system/server

# Usar venv de Python
Environment="PATH=/home/usuario/facial-attendance-system/server/venv/bin"

# Comando para iniciar Uvicorn
ExecStart=/home/usuario/facial-attendance-system/server/venv/bin/uvicorn \
    app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-level info

# Reiniciar automáticamente si falla
Restart=always
RestartSec=10

# Límites de recursos
LimitNOFILE=65536

# Logs a systemd journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
